# bebop_ws/devel/etc/catkin/profile.d/20.transform.bash
/home/cezar/bebop_ws/devel/.private/topic_tools/etc/catkin/profile.d/20.transform.bash